# overviewer/crontab
* * * * * /usr/bin/flock -n /tmp/ov.lock /usr/bin/overviewer.py --config=/config >> /var/log/cron.log
# this line is supposed to be here